# files/upstart.conf
# Emperor uWSGI script

description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

exec /usr/local/bin/uwsgi --ini /etc/uwsgi/emperor.ini --die-on-term
